# linux/ocsv003/.crontab
23 4 * * * sh ~/mytool/backup/cron.sh > ~/backup/cron.out 2> ~/backup/cron.err
