# linux/ocsv001/.crontab
41 2 * * * sh ~/mytool/backup/cron-host.sh > ~/backup/cron.out 2> ~/backup/cron.err
